# deploy/backend.service
[Unit]
Description=Lost Found Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/lost-found-system
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/lost-found-system/venv/bin/python backend/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
